# crontab
* * * * * root . /etc/environment; /app/run_moca_info.sh >> /var/log/cron.log 2>&1
