# crontab
* * * * * root /usr/bin/env bash -c 'cd /app && ./moca_info.py >> /var/log/cron.log 2>&1'
